# cron/crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
0 0 * * * certbot renew
* * * * * sh /hdd/backups/backup.sh >> /var/log/cron.log 2>&1
